# deployment_scripts/puppet/modules/telemetry/templates/hindsight/init.conf.erb
description "Telemetry collector hindsight service"

start on runlevel [2345]
stop on shutdown

exec start-stop-daemon --start -c hindsight --exec /usr/bin/hindsight -- <%= @conf_dir %>/hindsight.cfg
